# systemd/user/voicecall-ui-prestart.service.d/exec-silica-qt5-booster.conf
[Service]
ExecStart=
ExecStart=/usr/bin/invoker -o -s --type=silica-qt5 /usr/bin/voicecall-ui -prestart
