# systemd/user/voicecall-ui-prestart.service.d/silica-qt5-booster.conf
[Unit]
After=booster-silica-qt5.service
Requires=booster-silica-qt5.service

[Service]
ExecStart=
ExecStart=/usr/bin/invoker -o -s --type=silica-qt5 /usr/bin/voicecall-ui -prestart
